# package/lib/systemd/system/tlp-switch-daemon.service
[Unit]
Description=TLP-Switch Daemon
After=tlp.service

[Service]
ExecStart=sh -c "tlp-switch-daemon"

[Install]
WantedBy=graphical.target

